# system/services/search/author-updates-subscriber/fake-authors-service/supervisord.conf
[program:fake-authors-service]
command=java -classpath classes/:internal-libraries/*:external-libraries/* com.ssrn.search.author_updates_subscriber.fake_authors_service.Service server
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0