# system/services/frontend/website/supervisord.conf
[program:website]
command=java -classpath classes/:resources/:internal-libraries/*:external-libraries/* com.ssrn.frontend.website.Service server
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:filebeat]
autostart=%(ENV_SHIP_LOGS)s
command=/usr/share/filebeat/bin/filebeat -e -c /etc/filebeat/filebeat.yml
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0