# system/test/support/fake-old-platform/supervisord.conf
[program:fake-old-platform]
command=java -classpath classes/:resources/:internal-libraries/*:external-libraries/* com.ssrn.fake_old_platform.Service server
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0