# system/test/support/fake-kinesis/supervisord.conf
[program:kinesis]
command=kinesalite --ssl --deleteStreamMs 0 --createStreamMs 0 --updateStreamMs 0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:dynamodb]
command=java -Djava.library.path=./DynamoDBLocal_lib -jar DynamoDBLocal.jar -sharedDb -inMemory
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:configure-kinesis]
startsecs=0
environment=
    AWS_ACCESS_KEY_ID=dummy,
    AWS_SECRET_ACCESS_KEY=dummy
command=configure-kinesis.sh papers fake-papers-service-healthchecks fake-authors-service-healthchecks author-updates
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0